# earlyoom.service
[Service]
StandardOutput=null
StandardError=syslog
ExecStart=${PREFIX}/bin/earlyoom

[Install]
WantedBy=multi-user.target
