# earlyoom.service
[Service]
StandardOutput=null
StandardError=syslog
ExecStart=/usr/local/bin/earlyoom

[Install]
WantedBy=multi-user.target
