# earlyoom.service.in
[Unit]
Description=Early OOM Daemon
Documentation=man:earlyoom(1) https://github.com/rfjakob/earlyoom

[Service]
EnvironmentFile=-:SYSCONFDIR:/default/earlyoom
StandardOutput=null
StandardError=syslog
ExecStart=:TARGET:/earlyoom $EARLYOOM_ARGS

[Install]
WantedBy=multi-user.target
